# systemd/kimi-objective-heartbeat.timer
[Unit]
Description=Run Kimi objective heartbeat every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
Unit=kimi-objective-heartbeat.service

[Install]
WantedBy=timers.target
